# roles/04_application_geth/templates/geth.service.j2
# /etc/systemd/system/geth.service
# Hardened Geth service unit file

[Unit]
Description=Geth Ethereum Node (Secure)
After=network.target

[Service]
User=geth
Group=geth
Type=simple
ExecStart=/usr/bin/geth \
  --datadir /var/lib/geth \
  --http \
  --http.addr 0.0.0.0 \
  --http.port 8545 \
  --http.api eth,net,web3 \
  --authrpc.jwtsecret /var/lib/geth/jwt.secret \
  --syncmode "snap"

# Security hardening options
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ReadWritePaths=/var/lib/geth
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_CHOWN CAP_DAC_OVERRIDE
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target

