# roles/04_application_geth/templates/geth.service.j2
[Unit]
Description=Geth Ethereum Node
After=network.target

[Service]
User=geth
Group=geth
# Thay đổi ExecStart theo cài đặt thực tế (đường dẫn binary & options)
ExecStart=/usr/bin/geth --config /etc/geth/config.toml
Restart=always
RestartSec=10
LimitNOFILE=65536

# --- Sandbox security settings ---
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
ReadWritePaths=/var/lib/geth
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_CHOWN CAP_DAC_OVERRIDE

# Giới hạn _nproc_ và memory nếu muốn (tuỳ môi trường)
#TasksMax=100
#MemoryLimit=1G

[Install]
WantedBy=multi-user.target

